# crontab
*/30 * * * * /update-feed.sh
   0 3 * * 2 /update-wantlist.sh
