# crontab
*/30 * * * * /update-feed.sh
   0 3 * * 4 /update-wantlist.sh
