# crontab
*/30 * * * * /update-feed.sh
   0 3 * * 5 /update-wantlist.sh
